# supervisord.conf
[unix_http_server]
file=/kypost/state/supervisor.sock
chmod=0700
chown=kypost:kypost

[supervisord]
nodaemon=true
user=kypost
logfile=/kypost/logs/supervisord.log
pidfile=/kypost/state/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///kypost/state/supervisor.sock

[program:api]
command=/usr/local/bin/kypost-server --mode server
autostart=true
autorestart=true
startsecs=3
priority=10
stdout_logfile=/kypost/logs/api.log
stderr_logfile=/kypost/logs/api.err.log

[program:daemon]
command=/usr/local/bin/kypost-server --mode daemon
autostart=true
autorestart=true
startsecs=3
priority=10
stdout_logfile=/kypost/logs/daemon.log
stderr_logfile=/kypost/logs/daemon.err.log

[program:ollama]
command=/bin/bash -c "set -o pipefail; /bin/sh /opt/kypost/scripts/start-ollama.sh 2>&1 | while IFS= read -r line; do echo \"[$(date '+%%Y-%%m-%%d %%H:%%M:%%S')] $line\"; done"
environment=HOME="/home/kypost",OLLAMA_MODELS="/kypost/ollama-models"
autostart=true
autorestart=true
startsecs=3
priority=10
stdout_logfile=/kypost/logs/ollama.log
stderr_logfile=/kypost/logs/ollama.err.log

[program:ollama-model]
command=/bin/sh /opt/kypost/scripts/pull-ollama-model.sh
environment=HOME="/home/kypost",OLLAMA_MODELS="/kypost/ollama-models"
autostart=true
autorestart=false
startsecs=0
priority=11
stdout_logfile=/kypost/logs/ollama.log
stderr_logfile=/kypost/logs/ollama.err.log
